# systemd/meals.service
[Unit]
Description=Get meals from UA

[Service]
Type=oneshot
ExecStart=python ".ua-meals/meals.py"
